# photobooth.service
[Unit]
Description=Photobooth Service
After=network.target

[Service]
# setup_pi.sh replaces these placeholders during install.
User={{PHOTOBOOTH_USER}}
WorkingDirectory={{PHOTOBOOTH_DIR}}
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/{{PHOTOBOOTH_USER}}/.Xauthority

ExecStart={{PHOTOBOOTH_PYTHON}} {{PHOTOBOOTH_SCRIPT}}

# Restart automatically if it crashes
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
