# photobooth.service
[Unit]
Description=Photobooth Web Server
After=network.target

[Service]
# setup_pi.sh replaces these placeholders during install.
User={{PHOTOBOOTH_USER}}
WorkingDirectory={{PHOTOBOOTH_DIR}}
Environment=PYTHONUNBUFFERED=1

ExecStart={{PHOTOBOOTH_PYTHON}} server.py

# Restart automatically if it crashes
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
